# aget-state-tray.service
[Unit]
Description=AI inference server state tray indicator
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=%h/.local/bin/aget-state-tray
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
