# llama-server.service
[Unit]
Description=llama.cpp server
After=network.target

[Service]
EnvironmentFile=-%h/.config/llama-server/current.env
ExecStart=%h/llama-mtp/build/bin/llama-server --model ${LLAMA_MODEL} $LLAMA_ARGS
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
